# supervisord.conf
[supervisord]
nodaemon=true ; 关键！让 Supervisor 在前台运行，以保持容器存活

[program:python_app]
command=python run_camoufox.py config.yaml
directory=/app  ; 指定程序运行的目录
autostart=true
autorestart=true
stdout_logfile=/dev/stdout ; 将标准输出重定向到容器的 stdout
stdout_logfile_maxbytes=0  ; 禁用日志轮转
stderr_logfile=/dev/stderr ; 将标准错误重定向到容器的 stderr
stderr_logfile_maxbytes=0  ; 禁用日志轮转

[program:go_app]
command=/app/go_app_binary ; 这是我们编译后的 Go 程序路径
directory=/app ; Go 程序可以在 app 根目录运行，或指定其自己的目录
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0